# debian/ledgersmb.postinst
#!/bin/sh
# postinst script for ledgersmb
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

#  Set old_version variable for use later in the script.
old_version=$2

#   Remove these settings for /etc/ledgersmb/*, as dpkg-statoverride
# is no longer being used in debian/postinst to set it
drop_statoverride() {
    for file in /etc/ledgersmb /etc/ledgersmb/images \
        /etc/ledgersmb/images/demo /etc/ledgersmb/template ;
    do
        if dpkg-statoverride --list $file >/dev/null; then
            dpkg-statoverride --remove $file
        fi
    done
}

dbc_first_version="1.3.14-2"
dbc_dbtypes="pgsql"

. /usr/share/debconf/confmodule
<<<<<<< HEAD
=======
. /usr/share/dbconfig-common/dpkg/postinst.pgsql
>>>>>>> dbconfig-dev

case "$1" in
    configure)
        # Do Apache related processing if it is present.
        if [ -d "/etc/apache2" ]; then
            # Only do these if this is the initial install.
            if [ -z "$old_version" ]; then
                # Enable Apache mod_rewrite
                if [ -x "`which a2enmod 2>/dev/null`" ]; then
                    a2enmod rewrite
                fi
                # Restart apache
                if [ -x "/etc/init.d/apache2" ]; then
                    invoke-rc.d apache2 restart || true
                fi
            fi
        fi
        #   Remove these settings for /etc/ledgersmb/*, as dpkg-statoverride
        # is no longer being used in debian/postinst to set it.
        if [ -n "$old_version" ] &&  dpkg --compare-versions $old_version lt 1.3.9-2; then
            drop_statoverride
        fi
        # If enabled, do the debconf supported configuration
        db_get ledgersmb/lsmb_debconf
        if [ "$RET" = true ]; then 
            dbc_go ledgersmb $@
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
